# dot-config/systemd/user/zie-toolbox-image.timer
[Unit]
Description=image pulled every Mon, 01:05
[Timer]
OnCalendar=Mon, 01:05
AccuracySec=12h
Persistent=true
[Install]
WantedBy=timers.target

